# etc/supervisord.d/rabbitmq.ini
[program:rabbitmq]                                                                                          
command=/bin/sh -c '/usr/sbin/rabbitmq-plugins enable --offline rabbitmq_amqp1_0; /usr/sbin/rabbitmq-server'
